# trunk/debmarshal.debian.net/etc/cron.d/upload-foswiki
MAILTO=debmarshal@kernwifi.net
5,15,25,35,45,55 *    * * *   debmarshal    /usr/local/bin/upload-foswiki
