# debmarshal.debian.net/etc/cron.d/mirror-sid
MAILTO=dld@kernwifi.net
24 04    * * *   debmarshal    /usr/local/bin/mirror-sid
